# root/etc/services.d/eco-server/run
#!/bin/bash

su -s /bin/sh -c 'mono /app/eco/EcoServer.exe -nogui' abc
